# racksview_app-1.1/systemd/mdetector-front.service
[Unit]
Description=Ultrasonic motion detector on /dev/serial1
After=network.target

[Service]
ExecStartPre=/usr/bin/rm -f /var/log/racksview/mdetector-front.log
ExecStart=/opt/racksview/bin/mdetector-service.sh --profile front
Restart=on-failure
RestartSec=10
StandardOutput=file:/var/log/racksview/mdetector-front.log
StandardError=file:/var/log/racksview/mdetector-front.log
MemoryMax=50M

[Install]
WantedBy=multi-user.target
