# racksview_app-1.1/systemd/mdetector-back.service
[Unit]
Description=Ultrasonic motion detector on /dev/serial0
After=network.target

[Service]
ExecStartPre=/usr/bin/rm -f /var/log/racksview/mdetector-back.log
ExecStart=/opt/racksview/bin/mdetector-service.sh --profile back
Restart=on-failure
RestartSec=10
StandardOutput=file:/var/log/racksview/mdetector-back.log
StandardError=file:/var/log/racksview/mdetector-back.log
MemoryMax=50M

[Install]
WantedBy=multi-user.target
